# config/mytodo.conf
[program:my_todo]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.config.py
directory=/var/www/mytodo
autostart=true
autorestart=true
